# prosody-addon/intelli-block/cron.d/s2s_intelli_blacklist
#!/bin/bash
/usr/local/bin/update_s2s_blacklist.sh >> /var/log/s2s_blacklist_update.log 2>&1
